# claude/.config/systemd/user/claude-zulip.service
[Unit]
Description=Claude Code to Zulip Sync
After=network.target

[Service]
Type=simple
EnvironmentFile=%h/.secrets.env
ExecStart=%h/.claude/hooks/zulip-sync-launcher.sh
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
